# firefly/pes/eeb-temp.service
[Unit]
Description=Contiki-based temperature monitoring launcher
After=network.target

[Service]
ExecStart=/home/pi/contiki-3.0/contiki-projects/firefly/eeb/launch.sh
Restart=always
User=pi
Group=pi

[Install]
WantedBy=multi-user.target

